# setup/picar.cron
*/10 * * * * root /home/pi/setbrightness.ps1
